# OCP_Environment/suspend-prevent.service
[Unit]
Description=Prevent suspend on lid close when on AC

[Service]
Type=simple
ExecStart=/usr/local/bin/suspend-prevent.sh --forever
